# docker/rootfs/etc/s6-overlay/s6-rc.d/backend/finish
#!/command/execlineb -P
s6-svscanctl -t /run/service
