# docker/rootfs/etc/s6-overlay/s6-rc.d/backend/run
#!/command/execlineb -P
with-contenv
/bin/sh -c "
: \"\${ACCOUNT_ROLE:=all}\"
: \"\${JAVA_OPTS:=}\"
export ACCOUNT_HTTP_PORT=8081
exec java \$JAVA_OPTS -jar /app/app.jar serve --role=\$ACCOUNT_ROLE
"
