# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/run/supervisord.pid

[program:backend]
command=java %(ENV_JAVA_OPTS)s -jar /app/app.jar serve --role=%(ENV_ACCOUNT_ROLE)s
environment=ACCOUNT_HTTP_PORT="%(ENV_ACCOUNT_HTTP_PORT_INTERNAL)s"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=false
startretries=0

[program:nginx]
command=nginx -g "daemon off;"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=false
startretries=0

[eventlistener:exit_on_any_stop]
command=python3 /usr/local/bin/supervisor-event-listener.py
events=PROCESS_STATE_EXITED,PROCESS_STATE_FATAL
